# bdsd-io.service
[Unit]
Description=Socket.IO support for Bobaos Datapoint Sdk Daemon

[Service]
ExecStart=/usr/bin/env bdsd-io

[Install]
WantedBy=default.target
